# packaging/systemd/ikoma-mcp-runner.service
[Unit]
Description=IKOMA MCP Runner (passive)
After=network.target

[Service]
Type=simple
Restart=always
User=ikoma
Group=ikoma
WorkingDirectory=/opt/ikoma/code
EnvironmentFile=/etc/ikoma/runner.env
ExecStart=/opt/ikoma/venv/bin/python -m ikoma_mcp.runner.main

[Install]
WantedBy=multi-user.target
